# deploy/systemd/oecu-bot-update-dev.timer
[Unit]
Description=One-minute development update check for oecu Discord bot

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=oecu-bot-update.service

[Install]
WantedBy=timers.target
